# systemd/idle-shutdown.service
[Unit]
Description=Idle Auto-Shutdown (2 hours)
After=graphical.target
Wants=graphical.target

[Service]
Type=simple
ExecStart=/opt/ghost/scripts/idle_shutdown.sh
Restart=on-failure
RestartSec=60
Environment=DISPLAY=:0
Environment=HOME=/root

[Install]
WantedBy=graphical.target
