# systemd/idle-shutdown.service
[Unit]
Description=Idle Auto-Shutdown (2 hours)
After=graphical-session.target
Requires=graphical-session.target

[Service]
Type=simple
ExecStartPre=/bin/sleep 60
ExecStart=/opt/ghost/scripts/idle_shutdown.sh
Restart=on-failure
RestartSec=120

[Install]
WantedBy=graphical-session.target
